# deploy/jomove.service
[Unit]
Description=Jomove rental dashboard
After=network.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/jomove
ExecStart=/usr/local/bin/bun jomove.ts serve --port 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
