# deploy/jomove-scrape.service
[Unit]
Description=Jomove auto-scrape (portals → markdown → ingest → notify)
After=network-online.target jomove.service
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/jomove
ExecStart=/usr/local/bin/bun jomove.ts auto-scrape --no-archive
ExecStartPost=/usr/bin/touch /home/ubuntu/jomove/data/.last-auto-scrape
StandardOutput=journal
StandardError=journal
TimeoutStartSec=15min

[Install]
WantedBy=multi-user.target
